# scripts/depth-stream.service
[Unit]
Description=Depth Stream Server for G1 Robot (RealSense)
After=network.target

[Service]
Type=simple
User=unitree
ExecStart=/home/unitree/g1_inspector/build/depth_stream_server 192.168.123.255 5001 15
Restart=always
RestartSec=5
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
